# systemd/thermal-ctrl.service
[Unit]
Description=Thermal control simulation and validation harness
After=network.target
[Service]
Type=simple
ExecStart=/usr/bin/python3 -m thermal_ctrl dry-run --config /opt/thermal-ctrl-harness/configs/config.yaml
Restart=always
RestartSec=2
User=root
[Install]
WantedBy=multi-user.target
